# bin/ottoapp.service
[Unit]
Description=Otto App server
StartLimitIntervalSec=0
After=network-online.target

[Service]
Type=simple
User=ottoapp
Group=ottoapp
PIDFile=/run/ottoapp.pid
WorkingDirectory=/var/www/ottomap.mdhenderson.com
ExecStart=/var/www/ottomap.mdhenderson.com/bin/ottoapp serve --database /home/ottoapp/data/ottoapp.db
ExecReload=/bin/kill -USR1 $MAINPID
Restart=on-failure
RestartSec=13

[Install]
WantedBy=multi-user.target
